# PYTHON/REMOTE-PC-WEB-STATS-LINUX - docker (OG)/python-device-control-server.timer
[Unit]
Description=Timer for the python-stats-server service

[Timer]
OnBootSec=30sec

[Install]
WantedBy=timers.target
